# neo4j.ini
[program:neo4j]
command=/var/lib/neo4j/bin/neo4j console
